# etc/systemd/system/ptrrecordsync.timer
[Unit]
Description=Update PTR Records

[Timer]
OnBootSec=0min
OnUnitActiveSec=30min

[Install]
WantedBy=multi-user.target
